# update-extend-host.path
[Unit]
Description=Generate hosts file from leases file then reload dnsmasq

[Path]
PathModified=/var/lib/misc/dnsmasq.leases

[Install]
WantedBy=multi-user.target
